# config/systemd/user/bucket.service
[Unit]
Description=Backup Script

[Service]
Type=oneshot
ExecStart=/usr/bin/zsh -c 'bucket'

[Install]
WantedBy=default.target
