# config/systemd/user/bucket.timer
[Unit]
Description=Runs bucket.service periodically

[Timer]
OnCalendar=Fri 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
