# config/autostart/blueman.desktop
[Desktop Entry]
Name=Blueman Applet
Exec=blueman-applet
Type=Application
Hidden=true
