# config/autostart/nm-applet.desktop
[Desktop Entry]
Name=NetworkManager Applet
Exec=nm-applet
Type=Application
NoDisplay=true
Hidden=true
